# supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
command=/usr/local/bin/uwsgi --ini /anfisa/uwsgi.ini

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
